# conf/nginx.conf
[program:nginx]
command=nginx -c /var/lib/faucet/nginx.conf

stdout_logfile=/var/log/nginx.log
redirect_stderr=true
autorestart=true
autostart=true
